# templates/launchd/oss-repo-agent-dispatch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mikolaj92.hermes.repo-issue-to-pr-dispatch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mini-m4-main/.hermes/scripts/repo_issue_to_pr_dispatch.sh</string>
    <string>--live</string>
    <string>--run-opencode</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Background</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Keep unsafe Claude execution disabled unless explicit human approval and sandboxing review have happened. -->
    <key>HERMES_ALLOW_UNSAFE_CLAUDE</key>
    <string>0</string>
  </dict>
</dict>
</plist>
