# templates/launchd/oss-repo-agent-fala-tick-all.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Immutable Fala candidate template; renderer must replace every marker. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mikolaj92.hermes.repo-agent-fala-tick-all</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{UV_BIN}}</string>
    <string>run</string>
    <string>--frozen</string>
    <string>--project</string>
    <string>{{PROJECT_ROOT}}</string>
    <string>repo-agent-tick-all</string>
    <string>--config</string>
    <string>{{CONFIG_PATH}}</string>
    <string>--db</string>
    <string>{{DB_PATH}}</string>
    <string>{{MODE_ARG}}</string>
    <string>--json</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>{{PROJECT_ROOT}}</string>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>{{LOG_DIR}}/repo-agent-fala-tick-all.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{LOG_DIR}}/repo-agent-fala-tick-all.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
  <key>RunAtLoad</key>
  <false/>
  <key>LimitLoadToSessionType</key>
  <string>Background</string>
</dict>
</plist>
