# templates/launchd/oss-repo-agent-pr-triage.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mikolaj92.hermes.repo-pr-triage</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mini-m4-main/.hermes/scripts/repo_pr_triage.sh</string>
    <string>--live</string>
    <string>--comment</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Background</string>
</dict>
</plist>
