# launchd/com.hermes.oss-repo-agent.dispatch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template — replace {{HOME}} with your home directory before loading.        -->
<!-- Install: cp rendered.plist ~/Library/LaunchAgents/ && launchctl load -w ... -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.oss-repo-agent.dispatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{HOME}}/.hermes/scripts/repo_issue_to_pr_dispatch.sh</string>
    <string>--live</string>
    <string>--run-opencode</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Background</string>

  <key>WorkingDirectory</key>
  <string>{{HOME}}</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PATH</key>
    <string>{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HERMES_ISSUE_TO_PR_DRY_RUN</key>
    <string>0</string>
    <key>HERMES_ISSUE_TO_PR_RUN_OPENCODE</key>
    <string>1</string>
    <key>HERMES_ISSUE_TO_PR_BLOCK_INTAKE</key>
    <string>0</string>
    <!-- Max concurrent claude agents (default 3) -->
    <!-- <key>HERMES_ISSUE_TO_PR_MAX_CLAUDE_AGENTS</key> -->
    <!-- <string>3</string> -->
  </dict>

  <key>StandardOutPath</key>
  <string>{{HOME}}/.hermes/logs/repo-issue-to-pr-dispatch.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{HOME}}/.hermes/logs/repo-issue-to-pr-dispatch.launchd.err.log</string>
</dict>
</plist>
